# rc.local
#!/bin/bash
sh /data/etc/dbus-mqtt-services/reinstalllocal.sh
